# launchagent/local.chrome-log.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.chrome-log</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/modha/.claude/.venv/bin/python</string>
        <string>/Users/modha/Repos/skill-chrome-log/scripts/daemon.py</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/modha/.chrome-debug/logs/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/modha/.chrome-debug/logs/daemon.log</string>
</dict>
</plist>
